# srv/salt/prod/app/files/cmbcs.ini
[program:{{ appname }}]
command={{ command }}
process_name=%(program_name)s
directory={{ directory }}
user={{ user }}
autorestart=true
stdout_logfile={{ logs }}
